# files/polygram-network.service
[Unit]
Description=Polygram Network Initialization Service

[Service]
ExecStart=/usr/bin/polygram-init-network

[Install]
WantedBy=multi-user.target
